# influxdb/templates/crontab
0 5		* * *		root	/usr/local/sbin/influxdb_backup.sh {{ ansible_default_ipv4.address }} {{ influxdb_backup_dir }}
