# influxdb/templates/crontab
0 5		* * *		root	/usr/local/sbin/influxdb_backup.sh {{ influxdb_host }} {{ influxdb_backup_dir }}
